# ansible/templates/secret-manager.service.j2
[Unit]
Description=Secret Manager Service
After=network.target
Documentation=https://github.com/{{ secret_manager_repo }}

[Service]
Type=simple
User={{ secret_manager_service_user if secret_manager_service_user != 'root' else 'root' }}
Group={{ secret_manager_service_group if secret_manager_service_group else secret_manager_service_user }}
WorkingDirectory={{ secret_manager_work_dir }}
ExecStart={{ secret_manager_work_dir }}/secret-manager --config {{ secret_manager_config_dir }}/secret-manager.yml
Restart=always
RestartSec=5

# Security hardening
# NoNewPrivileges=true
# PrivateTmp=true
# ProtectSystem=strict
# ProtectHome=true
# ReadWritePaths={{ secret_manager_work_dir }} {{ secret_manager_config_dir }} {{ secret_manager_log_dir }} {{ secret_manager_data_dir }}

# Environment variables
# EnvironmentFile={{ secret_manager_config_dir }}/secret-manager.yml

[Install]
WantedBy=multi-user.target